# supervisord.conf
; ============================================================================
; Supervisor configuration for the Management Container
; ============================================================================
; Manages two primary processes:
;   1. webserver   — Python web server + API for managing Minecraft servers
;   2. velocity    — Velocity proxy for routing player connections
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ---- MinecraftServerTool Web Server ----
[program:webserver]
command=python webserver.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ---- Velocity Proxy ----
; Velocity starts after the webserver so DB init happens first.
; The velocity.toml is generated by the Python tool before/during startup.
[program:velocity]
command=java -Xmx512M -Xms512M -jar velocity.jar
directory=/app/data/velocity
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
